# systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=%h/bin/desktop/polybar-start.sh
ExecStop=killall -q polybar